# scripts/launchd/com.lifecoach.insights.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.lifecoach.insights</string>

    <!-- Runs `lifecoach insights generate` every day at 07:30 local time —
         after the daily reflection has run, so the insighter can use it. -->
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd "__REPO_ROOT__" &amp;&amp; /usr/local/bin/pnpm -w run lifecoach insights generate</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>7</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/data/logs/insights.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/data/logs/insights.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
