# scripts/launchd/com.lifecoach.daily-reflect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.lifecoach.daily-reflect</string>

    <!-- Runs `lifecoach reflect daily` every day at 06:00 local time.
         Update REPO_ROOT below to your actual repo path during install. -->
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd "__REPO_ROOT__" &amp;&amp; /usr/local/bin/pnpm -w run lifecoach reflect daily</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/data/logs/daily-reflect.out.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/data/logs/daily-reflect.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
